# dotfiles/dot_config/profile.d/00-profile.sh.tmpl
#!/bin/sh

{{ if eq .chezmoi.os "darwin" -}}
export PROFILE_OS=darwin
{{ else if eq .chezmoi.os "linux" -}}
export PROFILE_OS=linux
{{ else -}}
export PROFILE_OS=unknown
{{ end }}
