# dotfiles/private_dot_config/profile.d/20-git.sh
#!/bin/sh

export GITHUB_TOKEN="{{ .secrets.github.cli_token }}"
