# dotfiles/private_dot_config/profile.d/encrypted_05-work.sh.tmpl.age
-----BEGIN AGE ENCRYPTED FILE-----
YWdlLWVuY3J5cHRpb24ub3JnL3YxCi0+IFgyNTUxOSA2MFJaamcrVHhYY3UzQ2xj
dS9BK3MrYUhyNEpJQkJTZS9Qc1BBTjNvcW5rClRud090bHpPSkhTMnhEWWF5NCta
TllOenpUcDJ6OU12L1JWSWxzL3BzbE0KLS0tIFJlMVhEb3pXQVdjOWdRSm1pQ1hI
a3ZzVUZwZmE3dUNUZUNjcTZXa0lsbHMKABHlgME3OqdOLN0EnM8h24AMT4zJsyRy
wgKC7z59hOuACHKxneQ4aw8FOqHzXYrGPynI526TlLnXnnqeGY0FIPukXc4WACU4
mCEuUWbZgZm1vdV9cQTdc3DVe3Ezs5LwAV1z84gqRm18yn+ETkuiTLsyXkDr0iJn
UQJT1QDKBMDbmdb7m92+AyRF0SKWPs+i5MHrW/Wnt2Rb/UJO/tRCfw+syXs50+TQ
JD1/HRHvqjXnFpRN0ZiiDNH5UOotd7CSopQeGp6e7aWnx01tPyS9n4nz9CyDj+Yv
TKPMCkL8T73oWSNGWovNJbraP35ClSK0402+QG6V6UVHSRAw2Bh91nob2jmuFTYf
h0RxpoWzyr555IzISE9AxoqTRQprZO6cbZV8tnv+WUtXhuNUnyJKZlEGTbKaAviD
QvlbTJ0gnSUiwQsfyvHvE5sT1+/XfoUKucaduD50lvUpTWG8D9lhNJAbHpwqr7F9
NSGnpmm2S3ya3bQbcVlsXSwZ03FweDDJSKE0y1iYFg5lTp2odPCSieZjSqydWMal
UEdOi59pycs6gs02l/KPiM2ao3HakY1sFSZENu9fkleBTvUTV4BAv6+wfOETqrX1
4eu2wNUcI13ECfMTjSreiO+E9zUWiWAKWbhgOF9Ss5BplHOKZdV6mrtyniv90XwD
+P+tIGGzJ2kcmVdiWCDFd5GJo+LyfItZKhgMi/v7jyqSVLR+MbLcekeEGWglomfp
tPsekLzWyNDqHTjXL3PR4ZKTVxVi2oXEW/5qe13SDj7UkeLBKNrWi+z9mzLtfgpI
wJPkvzl/xcbQvdPiWfpKe1mE9B9l+ipE7DeN5GjOq8GMxS8BJkmzEkiqecwkEaco
APhhSr88vWf1yDej/mgkvvCOpql66X+FQQf0qygjW8jylPZHtS8P5VI/fVzya3Z/
hJ/lmOJn0K511xbEfGzgz+4VECpsGlXwgWf8BsfTK4QJ9ye2O2WC909ZcbRZ7zLU
7sb8EaZCMazAQN/t6hWv+M3v8TVBHU7NQfbl9pyhYLPvzYuy4fdTfEJIXERscv87
lVItSqCW5DaDi+PWSH/mQ2+oA2sWB8J56/nIJoYNmJ/uZSrrSNnMGdPw6D3IpT7K
oS4SRV/lupNpAZd5oc51pnglLdxstfbYf/DDWQ3sCEx01s9ikyNyb0DsrdjTxtXh
8TAvT+sF34JIq2KLwAS4Oia8fik+Jn9wRa7uUxhMWnfVXF5DAE2qn/y/VFYng+Y5
Oxnzlygakwp1rMWlqIbBg39uAH9tUVDNU55bgmEGIKnzWRa48QCDAedY88Y9PUJg
okX2FQ1ap1Vpmva0e0629re6DbpL89vWJ8WbP9GusgvkTnZ1CM/YjGLI8fJIRruB
tdXc14Tns8PLoAuMeMzJgl+ovsMa8MJxFS4si/yBNUnZyHZg9Y8jmvCmbh3GZ8Fv
iOZfQ80a2bj5o3g0qMMBKaU58xm8pRc8PGLdVim8odbDs05vUeQYXgRXtEBDeBvb
PYCP6mV51mslPsyaRUKcD9ipBlxcceqb
-----END AGE ENCRYPTED FILE-----
