# dotfiles/private_dot_config/profile.d/encrypted_05-work.sh.tmpl.age
-----BEGIN AGE ENCRYPTED FILE-----
YWdlLWVuY3J5cHRpb24ub3JnL3YxCi0+IFgyNTUxOSBrMEw0TE1aMENNS0xUTm1t
MjZGc1ZNRFgrZS9CZGRGVFRiZ2R5Z1pUSldVCkswWkxxL2U3YytsdTJISTRMdG5o
WDFPVUZRZEkvdEhtN2t1WGRnTG9hUFEKLS0tIHNXUys0VWUra2hyd2Q5S1dvOFNv
N3ZidkUyT0s1Y2VWRDBlcm1NaDJEM1EK0v6i1v6hAka72Ju0zC8ja2zvD2axGsrr
A4j6bA8zQ8F2/6K+4I4O+gMG5y33s51R+CS0uSoRy43W4XclGaCrRVMSkG/WI16Z
EXs5Y33TBK+ZUnqogAJpmb7SL0o6b5sJJQo8WfW60BWtAUKmYy5n5LOlcsIyncpI
QE7Ej9F7tEXNHTNhGYVyI3sMtLkbSVOMHCUhnFjA5H/Qv6C3oDhj3RC8cL1D9RoT
bB033hAgZkduA9Eyl6gYAfh3Uz2FCVjtlXxPuKcMTaBFbPhqdYnmPy/E/0NfYPCP
DaPLk8AUx0r7BE2SI0LR6ChPqPO06yjk3OHPcx1SEFhN92nu1Xh+pcE24tWvSOia
va/RLkCzVt06LfbAn6G7vBPVrxVUMqo3q++RCDULYYCb8pOQJlU51EnkBGdmIl5J
lk/cTkJn1N1NJbAttAPz4emg8ZFQHkZp5XJ0S1SqRZLe5uORq3uGiyu43trhhMTY
z8S4Z8NB1RYK7RMx9w3DreYrY1C8SObniixnyYZrubDFg8LYsIgX0i8Wedk+SAul
1BxljQfZ5FSAv1jXLP3Cr6yLudaHrgkuSz7oBYyHXSYVzCAOah7m0/QasSsu6qfH
mC+/fBhSBmCz5w0EExJnoJa3ooYV744x6O8EKzukcQ9P4/oo15GFjwKmbrTgfYjj
zsr7UN/NYo67FkfBq80dYbzvcCOXsie+snCnSabf1OYMRb/T0LunT4kbLnWNJz8P
qfhVlqmHI0NOZaBFrZ4EHea4aFre00l7K1ZQcB07Y3qFStfaqAJwgaWu1bcaNnJd
yhsFBosDpxTrKkNQMaXKK3cSrrIJQs2EGBJWy1YgXv6+eaJhZk3cF69IWnXNpRTn
aYQZPVqzorAiZsBrEc3cC+nzZ5woPzmuQApHsg97gW4CWKxzOPwSmhNfmanOiujC
uXmIjhugUYGnQzBl0fwsu6XuqKDGY5HrKEcr1FiOjr8WpU7A5srIBq2Wrx4MGHvb
xYWcMtEM5VijN3f7px4S91+b8rys91Kv2qwxCAXhWKhTGM7cKcgNdF5+YOlk+BtQ
PEvMoJWDU84+Tw1L94ffAzEtUGaBVnOq6KeKD64gm/tfV56km6zEK736cZ1EKb0P
RCOTQoET2kE9RZvcZwk+ZoBcpHdFwNdmoFN6vcJ4yX61wEI7e5AnExGE88PTnduV
w8FaVpRIzBAilBcMTuiCFOqTSPEFM4vZ+pBwrhaHBCEePJM+MZnYXkEW+V9s9T2T
JeQzGXMOFrYjZhkWvSC/z5Yn4Pqrzxu5De/NgXSFKyIl138fMBgceIA8NS3Ubrf7
UsDf7UyzQhpQj2dpKuKle1D6GZupd/cOAt1H+oWL5zlnjy74DDgQ5/0Pj5PuzEOD
KgsyW34tWuIlva0azE/uMvBYTQOM+sde2pSEUpOcZ0HvX5YnZejO8p2Ssci4YNlD
0AqbgajC93pFH7hcwB/f4rdUExDSETUECDLVPJitTy0ynXrlqf3FYggDZ6xmZsCt
khiNcr7TGMnrUB53dHuHOr3utOUGibZI7OY3eW/ww73cLUbhM8dAZNEOyj6CdDcl
I/KRyLTbg3wTJvItrq/bxQiR/Aot
-----END AGE ENCRYPTED FILE-----
