# dotfiles/private_dot_config/profile.d/encrypted_05-work.sh.tmpl.age
-----BEGIN AGE ENCRYPTED FILE-----
YWdlLWVuY3J5cHRpb24ub3JnL3YxCi0+IFgyNTUxOSAvOTRKeUpkMm0xcXBneHk4
YlBQMWNXWjhEMlVDUW8yMktxR09QdlU5MVZzCm5CYmRsUUU4UFdGWjExTWxUS1Ny
RVZqWTlrTzN0ZlNkL1dHL3h3SnlDTVkKLS0tIE1mbS80Z1BocVA2S3Avc3dzT3Iy
Qzdmb3lBcEtsdTRoc3dLQ2IvVmNlb0EKnd6lzZ0lGzlwOltHkoqDGTEJzhwF/1Ot
c+YC0H4PmnFxOycQBlGZWTwEltYUjCUKzqcIDAu8oG5P2X8RqojhpXWyQuMKrXJ9
iUZbNH7h3yeFjf8LuAmZNr5fGe1GfI85nKq88KeY13GvBiHb9ppAqF6RYKW0fMje
iiMSBl1g+c8pYTjtnYxmuELtSrUzREv+9SCOrCcit/I/5Bqk11pgV60iMAOPC/J8
DkPnrpxS5BcDU0bWRxq5rOovfoe0pES2NTmtSlppvYwoFC760Ltl2zQmXP0CwfpX
sr0CC6iNcYpeY2KV2LKd7A2Kc2bn0cGoc8Nv+kEEKVPcEbnQvB1w2c9nqAqFsCS5
Wxk4yMZ0umWKQpu0ybK0gU87y6y38pbhR1CWQFzdcgsy3C5hAWhMZu2IgdsCesCG
S+tdY79xpp72RLfCabl5UAWcNNtyTJMPsm9+9r3oSJ8TElRNryCI5NrTqA0IBpXI
e7t/OqS0kXHsYVmF7VBrE5C97UtHgiq/8YLxZZRBLmbbzbTcKISLJBDaWD5/luST
9UfKHllxUMKujJbJxqo3fjVHX9PNugeRZPdUz+h+hAMe2f0PbQYhhIwS6sOTLgQ4
tXvc9ZjQbODEOpBH1KdbXns8rUPvTT6yqL0HxdF+rDZYhJuKrzpn/MjLtFbLi92G
WYYrcxnWxsZUCJj3jyoPgA8rksuP0JSnSKZ2ZGGfytDdnR6qCuUgA4g5lVMx0DG0
zLkHAzMNwKiVnRnrXG5vag0/OaiOl7qBotP319simZGm+nTu6NaGDPeBr0kpBAdq
t2ipEjxIHQK3zbU9X18Blx7OFk3KJm4FWl/vdK4NiMiRhBT3v3V8rQal62X76xGy
FZRIPq/FiTCMdfrGRcmzDSSu3BKE6EHnCuSk9zoXunBiVRqDb3ibKOo1iAQ4Oigo
QRr1Xk52kc+m12REeutBDDypFUpcw9bhQnbLiTnR/JoQakbwvxYI750pB+gPcgN5
ABuCsjxVp9BbKLSfZPW3YFoyJXw08QFAJ8s/4XmIcIcdpVX2RcqBm+fQuRWCReDr
nEjKYVKzmMFBj6lAY7KWexjnhNHc3MwBYlrCLa2yDZpRwly1yPI1aIMLojeYsAPD
BF+65RbEXbse4dze+NKGByHtLdI8OGxUs3K0LaHvdVjRwUIwqzV3kaZF8oifqj1E
uLCu/rcVXvx1S4ej/zq4XT/FSm8hcrM2GSWlzugsE4VHLJ8EGEB8SmBnQeAGSEj9
mKWv7dmpvvkNSbFf0xF+VGd9dBqqRS88gHdsz8m7InqcA85zgNX3c1Mkdw5gmjz8
piehQ+tuqBjaR60gMEH5/C7yAvzzVQR1AM0IR9EFgQBUBPa+5oZ4dZZPOLOr8mF/
BM4cs5uTu9emTq5iMHY8cssoeRRwj6OiIlOCvM7qV/Qp+SMS7qJCm2YRz67XpZVR
JeB7NVdfYzrIPNY14CBmLg==
-----END AGE ENCRYPTED FILE-----
